# lky-admin.service
[Unit]
Description=LkyPanel — Admin Interface (port 2087)
After=network.target

[Service]
Type=simple
User=lkypanel
Group=lkypanel
WorkingDirectory=/usr/local/lkypanel
EnvironmentFile=/usr/local/lkypanel/.env
ExecStart=/usr/local/lkypanel/venv/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:2087 \
    --certfile /usr/local/lkypanel/ssl/panel.crt \
    --keyfile  /usr/local/lkypanel/ssl/panel.key \
    --access-logfile /var/log/lkypanel/admin-access.log \
    --error-logfile  /var/log/lkypanel/admin-error.log \
    --timeout 60 \
    lkypanel.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
